# poll-bme280.service
[Unit]
Description=Poll Service

[Service]
ExecStart=/opt/poll-bme280/poll-bme280.sh

[Install]
WantedBy=multi-user.target
